# gpio-monitor.service
[Unit]
Description=GPIO Monitor SSE Server
After=network.target
Documentation=https://github.com/WIPtitle/gpio-monitor

[Service]
Type=simple
User=pi
Group=gpio
ExecStart=/usr/bin/python3 /usr/lib/gpio-monitor/gpio-monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target