# gpio-monitor.service
[Unit]
Description=GPIO Monitor SSE Server with REST API
After=network.target
Documentation=https://github.com/WIPtitle/gpio-monitor

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/lib/gpio-monitor/gpio-monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/gpio-monitor

[Install]
WantedBy=multi-user.target